# teruez.service
[Unit]
Description=Teruez Web Server

[Service]
ExecStart=/usr/local/bin/teruez

[Install]
WantedBy=multi-user.target
